# cron.d/tasks-cron
# Minute Hour Day of Month Month Day of Week Command
# (0-59)(0-23) (1-31)(1-12 or Jan-Dec)(0-6 or Sun-Sat)

# Chosing shell and loading .env variables
SHELL=/bin/bash
BASH_ENV=/container.env

# Servidor nos EUA está 3 horas na frente do Brasil
30 21 * * * root /usr/local/bin/python /src/cron_calc.py >> /var/log/cron.log 2>&1 && root /usr/local/bin/python /src/bot.py >> /var/log/cron.log 2>&1

# Don't remove the empty line at the end of this file. It is required to run the cron job
